# launchd/com.martyuk.telegram-radar.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.martyuk.telegram-radar.daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>-i</string>
    <string>PATH=/usr/bin:/bin:/usr/sbin:/sbin</string>
    <string>HOME=/Users/martyuk</string>
    <string>USER=martyuk</string>
    <string>LOGNAME=martyuk</string>
    <string>/bin/zsh</string>
    <string>/Users/martyuk/Documents/telegram-radar/scripts/daily_collect_and_push.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/martyuk/Documents/telegram-radar</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>21</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/martyuk/Documents/telegram-radar/logs/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/martyuk/Documents/telegram-radar/logs/launchd.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
